# rules.d/84-kbdbacklight.rules
KERNEL=="smc::kbd_backlight", SUBSYSTEM=="leds", RUN+="/usr/bin/chmod 666 /sys/class/leds/%k/brightness"
