# rules.d/85-screenbacklight.rules
KERNEL=="nvidia_backlight", SUBSYSTEM=="backlight", RUN+="/usr/bin/chmod 666 /sys/class/backlight/%k/brightness"
